# deploy/hetzner/systemd/hapay-alert.timer
[Unit]
Description=Хапай — перевірка збору кожні 15 хв

[Timer]
# Кожні 15 хв. Поріг тиші — 45 хв (qtasks.COLLECT_SILENT_MIN), тож зупинку
# помітимо максимум за 60 хв, а не «колись уранці».
# Інтервал таймера має лишатись помітно меншим за поріг: інакше тривога спізнюється
# рівно на нього. 15 хв при порозі 45 дають три перевірки у вікні.
OnBootSec=10min
OnUnitActiveSec=15min
# Persistent — щоб після перезавантаження сервера перевірка не чекала повний інтервал
Persistent=true

[Install]
WantedBy=timers.target
